# contrib/debian/keep-web.service
[Unit]
Description=Keep Web Admin and FROST Co-Signer
Documentation=https://github.com/privkeyio/keep/tree/main/keep-web
After=network-online.target
Wants=network-online.target

[Service]
Type=exec
User=keep-web
Group=keep-web
ExecStart=/usr/bin/keep-web
Restart=on-failure
RestartSec=10
# A co-signer that stops retrying drops the group below threshold silently,
# so keep restarting rather than giving up after the default burst limit.
StartLimitIntervalSec=0
EnvironmentFile=/etc/keep-web/keep-web.env

LimitNOFILE=8192
TasksMax=512

StateDirectory=keep-web
StateDirectoryMode=0700
WorkingDirectory=/var/lib/keep-web
UMask=0077

NoNewPrivileges=yes
CapabilityBoundingSet=
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectProc=invisible
ProcSubset=pid
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
# AF_NETLINK is used by glibc getaddrinfo for source address selection and
# AF_UNIX by the nss-resolve/nscd lookup paths; without them name resolution
# degrades silently rather than failing.
RestrictAddressFamilies=AF_INET AF_INET6 AF_NETLINK AF_UNIX
MemoryDenyWriteExecute=yes
LockPersonality=yes
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

[Install]
WantedBy=multi-user.target
